# targets/firefox/config/^etc/systemd/system/firefox-suspend.service
[Unit]
Description=Freeze Firefox around suspend/resume
Before=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=-/usr/bin/pkill -STOP -x firefox
ExecStop=-/usr/bin/pkill -CONT -x firefox

[Install]
WantedBy=sleep.target
